# docker/cronjob/crontab
0 * * * * root /usr/local/bin/python3 /code/src/manage.py runcrons >> /var/log/cron.log 2>&1
