# docker/cronjob/crontab
*/1 * * * * python3 /code/src/manage.py runcrons
